# systemd/.config/systemd/user/battery-notify.service
[Unit]
Description=Battery notification service

[Service]
ExecStart=%h/.local/bin/battery_alert.sh
